# recipes-core/logrotate/files/logrotate-frequent.cron
*/10 * * * * root /usr/sbin/logrotate -v -l /var/log/logrotate.log /etc/logrotate.conf
